# src/pandastix_docker_supervisor.conf
[program:mysqld]
command=/pandastix_docker_mysql.sh
numprocs=1
autostart=true
autorestart=true

[program:apache2]
command=/pandastix_docker_apache.sh
numprocs=1
autostart=true
autorestart=true
